# data/unity8-policy-kit.conf.in

start on started unity8

exec @CMAKE_INSTALL_FULL_PKGLIBEXECDIR@/policykit-agent
